# Final deployment & Rest api/guni.conf
[program:gunicorn]
directory=/home/ubuntu/django
command =/home/ubuntu/env/bin/gunicorn --workers 3 --bind unix:/home/ubuntu/django/app.sock  happyvending.wsgi:application
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn/gunicorn.err.log
stdout_logfile=/var/log/gunicorn/gunicorn.err.log

[group:guni1]
programs:gunicorn
